# deploy/systemd/netcut-frontend.service
[Unit]
Description=Netcut frontend service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netcut
Group=netcut
WorkingDirectory=/home/netcut/frontend
EnvironmentFile=-/home/netcut/frontend/netcut-frontend.env
ExecStart=/bin/sh -lc 'exec ${NETCUT_FRONTEND_COMMAND:-/usr/bin/npx serve -s /home/netcut/frontend/dist -l 3000}'
Restart=always
RestartSec=5s
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/netcut/frontend
RuntimeDirectory=netcut-frontend
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
